# service/log/run
#!/bin/sh
exec setuidgid kafkafeeder multilog t n20 s10240000 /www/kafkafeeder/self-logs/
